# tools/release/systemd/pqnas-fail.service
[Unit]
Description=PQ-NAS rollback test service (always fails)
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/false
