# tools/release/systemd/pqnas-ok.service
[Unit]
Description=PQ-NAS swap test service (always succeeds)

[Service]
Type=oneshot
ExecStart=/bin/true
